# config/linux/systemd/timesketch.socket
[Unit]
Description=Timesketch socket

[Socket]
ListenStream=/run/timesketch/socket

[Install]
WantedBy=sockets.target
